# etc/cron.daily/openattic
#!/bin/sh

set -e

# Clean out expired Django sessions.
oaconfig clearsessions
