# debian/cron.d/updatetwraid
PATH=/usr/sbin:/usr/bin:/sbin:/bin

*/5 * * * * root cd /usr/share/openattic && ./manage.py updatetwraid -v0
