# debian/openattic-module-tftp.postinst
#!/bin/sh

set -e

if [ "$1" = "configure" ] ; then
    update-rc.d tftpd-hpa remove
fi
