# debian/openattic-module-tftp.postinst
#!/bin/sh

set -e

if [ "$1" = "configure" ] ; then
    invoke-rc.d tftpd-hpa stop
    update-rc.d -f tftpd-hpa remove
fi

#DEBHELPER#
