# debian/openattic-module-ftp-pgsql.postinst
#!/bin/bash

set -e

if [ "$1" = "configure" ] ; then
    dbconfig-generate-include -f template \
        -o template_infile=/etc/openattic/proftpd/pgsql.conf \
        /etc/dbconfig-common/openattic-pgsql.conf \
        /etc/proftpd/conf.d/openattic.conf

    if [ -x /etc/init.d/proftpd ]; then
        invoke-rc.d --quiet proftpd reload
    fi
fi


#DEBHELPER#
