# debian/openattic-module-tftp.postrm
#!/bin/bash

set -e

if [ "$1" = "purge" ] ; then
    rm -f /etc/xinetd.d/openattic-tftp
    invoke-rc.d xinetd reload
    update-rc.d tftpd-hpa defaults > /dev/null
    invoke-rc.d tftpd-hpa start
fi

#DEBHELPER#
