# debian/openattic-gui.postinst
#!/bin/sh

set -e

if [ "$1" = "configure" ] ; then
    if [ -x /etc/init.d/apache2 ]; then
        invoke-rc.d --quiet apache2 reload
    fi
fi

#DEBHELPER#
