# debian/openattic-base.postrm
#!/bin/bash

set -e

if [ "$1" = "purge" ] ; then
    update-rc.d openattic_systemd remove >/dev/null
    update-rc.d openattic_rpcd    remove >/dev/null

    rm -f /var/log/openattic_rpcd
    rm -f /var/log/openattic_systemd
    rm -f /var/log/openattic_vgmanager
    rm -rf /var/log/openattic
fi

#DEBHELPER#
