# etc/profile.d/openattic.csh
setenv PYTHONPATH "${PYTHONPATH}:/usr/share/openattic/nagios/conf"
