# etc/profile.d/openattic.sh
PYTHONPATH=$PYTHONPATH:/usr/share/openattic/nagios/conf
export PYTHONPATH
